# templates/etc/init.d/bucky2.Debian.erb
#!/bin/sh
#
# /etc/init.d/bucky2
# init script for Bucky2.
#
# Modified for LSB compliance by Glenn Aaldering <gaaldering@ebay.com>.

### BEGIN INIT INFO
# Provides:          bucky2
# Required-Start:    $syslog $named $network $time
# Required-Stop:     $syslog $named $network $time
# Should-Start:      
# Should-Stop:       
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Bucky2
# Description:       Bucky2
### END INIT INFO

NAME=bucky2
DESC="Bucky2"
DEFAULT=/etc/default/$NAME

if [ `id -u` -ne 0 ]; then
   echo "You need root privileges to run this script"
   exit 1
fi

. /lib/lsb/init-functions

if [ -r /etc/default/rcS ]; then
   . /etc/default/rcS
fi


if [ -f "$DEFAULT" ]; then
   . "$DEFAULT"
fi

PIDFILE='/var/run/bucky2.pid'
DAEMON='/usr/local/bin/bucky2'
CONFIG='/etc/bucky2/bucky2.conf'


RETVAL=0

bucky2_start() {
  log_daemon_msg "Starting $DESC"
  if pidofproc -p "$PIDFILE" >/dev/null; then
    log_action_end_msg 3 "Already running"
    exit 0
  fi
  start-stop-daemon --exec $DAEMON -b -S -m -p $PIDFILE -c <%= scope.lookupvar('bucky::run_as_user') %> -- $CONFIG
  log_end_msg $?
}

bucky2_stop() {
  log_daemon_msg "Stopping $DESC"
  start-stop-daemon -K -p $PIDFILE
  rm -rf $PIDFILE
  log_end_msg $?
}

bucky2_status() {
  if pidofproc -p "$PIDFILE" >/dev/null; then
    log_action_end_msg 0 "running"
    exit 0
  else
    if [ -e "$PIDFILE" ]; then
      log_action_end_msg 1 "failed to start"
      exit 1
    else
      log_action_end_msg 0 "not running"
      exit 3
    fi
  fi
}

case "$1" in
  start)
    bucky2_start
    ;;
  stop)
    bucky2_stop
    ;;
  restart)
    bucky2_stop
    bucky2_start
    ;;
  status)
    bucky2_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac
